# supervisord.conf
[supervisord]
nodaemon=true

[program:fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize --force-stderr

[program:nginx]
command=/usr/sbin/nginx
